# cookbooks/openruko/templates/default/upstart-openruko-httprouting.conf.erb
start on starting openruko
stop on stopping openruko
respawn

export HTTPROUTING_PORT=80; \
exec su - rukosan -c 'cd /home/rukosan/openruko/httprouting; \
export PG_USER="rukosan"; \
export PG_PASSWORD="rukosan"; \
export HTTPROUTING_TLS=true; \
httprouting/bin/httprouting >> /var/log/openruko/httprouting.log 2>&1'
